# packages/cloud-builder/cloud_builder-1.0.7.tar.gz/cloud_builder-1.0.7/systemd/cb-fetch-once.service
[Unit]
Description=CB Project Fetch Once
After=syslog.target
Requires=network-online.target

[Service]
EnvironmentFile=-/etc/cloud_builder
Type=simple
ExecStart=/usr/bin/cb-fetch --project "${CB_PROJECT}" --branch "${CB_BRANCH}" --single-shot

[Install]
WantedBy=multi-user.target
