# packages/cloud-builder/cloud_builder-1.0.7.tar.gz/cloud_builder-1.0.7/systemd/cb-info.service
[Unit]
Description=CB Project Info Service
After=syslog.target
Requires=network-online.target cb-fetch-once

[Service]
EnvironmentFile=-/etc/cloud_builder
Type=simple
Restart=on-failure
ExecStart=/usr/bin/cb-info --update-interval "${CB_UPDATE}" "${CB_INFO_RESPONSE_TYPE}"

[Install]
WantedBy=multi-user.target
